# scripts/systemd/agentdock.service
[Unit]
Description=AgentDock local daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/projects/ai-task-workbench
Environment=NODE_ENV=production
Environment=AGENTDOCK_SERVE_STATIC=1
Environment=PATH=%h/.nvm/versions/node/v24.13.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env npm run server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
